# templates/cron-daemon/systemd.service
# Genesis-Core Cron Daemon — systemd user unit (Linux)
#
# Installation:
#   1. Copy this file to ~/.config/systemd/user/genesis-cron-daemon.service
#   2. Edit ExecStart below to point to the genesis-core binary.
#   3. Optionally override GENESIS_HOME if it differs from ~/.genesis.
#   4. Enable + start:
#        systemctl --user daemon-reload
#        systemctl --user enable --now genesis-cron-daemon.service
#   5. Verify:
#        systemctl --user status genesis-cron-daemon.service
#
# Uninstall:
#   systemctl --user disable --now genesis-cron-daemon.service
#   rm ~/.config/systemd/user/genesis-cron-daemon.service
#   systemctl --user daemon-reload
#
# Logs:
#   journalctl --user -u genesis-cron-daemon.service -f
#   # or from the log file:
#   tail -f ~/.genesis/cron-daemon.log

[Unit]
Description=Genesis-Core Cron Daemon
Documentation=https://github.com/dmercer290-byte/wayland-core
After=network.target

[Service]
Type=simple

# Replace with the absolute path to the genesis-core binary.
# Example: ExecStart=/usr/local/bin/genesis-core cron daemon
ExecStart=/usr/local/bin/genesis-core cron daemon

# Run the daemon body directly (systemd is the process manager).
Environment=GENESIS_CRON_DAEMON_CHILD=1

# Override if your Genesis home is not ~/.genesis.
# Environment=GENESIS_HOME=/home/youruser/.genesis

# Redirect logs to the daemon log file. systemd also captures
# stdout/stderr via journald by default.
StandardOutput=append:%h/.genesis/cron-daemon.log
StandardError=append:%h/.genesis/cron-daemon.log

# Restart on unexpected exit. Do not restart on SIGTERM (clean shutdown).
Restart=on-failure
RestartSec=10s

# Allow the unit to start at login via `loginctl enable-linger`.
KillMode=process
TimeoutStopSec=10

[Install]
WantedBy=default.target
